# deploy/oracle/tempo-explorer.service
[Unit]
Description=INS TEMPO Explorer
After=network.target

[Service]
Type=simple
User=tempo
WorkingDirectory=/opt/tempo
Environment=TEMPO_DATA_DIR=/opt/tempo/data
Environment=TEMPO_DEBUG=false
ExecStart=/opt/tempo/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8080 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
